# systemd/heraldd.service
[Unit]
Description=Herald - Claude Code Telegram Remote Control Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=%h/.local/bin/heraldd
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
WatchdogSec=30
Environment=RUST_LOG=info

# Systemd creates these dirs before namespace setup, avoiding NAMESPACE errors
RuntimeDirectory=herald
ConfigurationDirectory=herald

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=%h/.config/herald %t/herald
PrivateTmp=yes

[Install]
WantedBy=default.target
